# files/launch/powermgr.conf.in
# @@@LICENSE
#
#      Copyright (c) 2015 LG Electronics, Inc.
#
# LICENSE@@@

description "@WEBOS_PROJECT_SUMMARY@"

start on minimal-boot-done
stop on started start_update

respawn

# Comment this line out to suppress logs on the console
#console output

pre-start script
@WEBOS_INSTALL_SBINDIR@/msm_suspend
@WEBOS_INSTALL_SBINDIR@/set_cpufreq_ondemand_governor
@WEBOS_INSTALL_SBINDIR@/cpufreq_governor_parameters
end script

exec @WEBOS_INSTALL_SBINDIR@/powermgr
