# files/launch/powermgr.conf.in
# @@@LICENSE
#
#      Copyright (c) 2017 LG Electronics, Inc.
#
# LICENSE@@@

description "@WEBOS_PROJECT_SUMMARY@"

start on minimal-boot-done
stop on started start_update

respawn

# Comment this line out to suppress logs on the console
#console output

exec @WEBOS_INSTALL_SBINDIR@/powermgr
